# app-emulation/qubes-core-qubesdb/files/qubesdb-daemon
#!/sbin/openrc-run
# $Id$

depend() {

	need qubes-selinux
	need tmpfiles.setup
}

start() {

	ebegin 'Starting qubesdb-daemon'

	# start-stop-daemon hangs =(
	#
	# start-stop-daemon --env --start --exec /usr/sbin/qubesdb-daemon -- 0

	/usr/sbin/qubesdb-daemon 0
}


stop() {

	ebegin 'Stopping qubesdb-daemon'

#	start-stop-daemon --stop --pidfile /run/qubes/qubesdb.pid

	pkill -15 qubesdb-daemon
}
