# app-emulation/qubes-core-qubesdb/files/qubesdb-daemon
#!/sbin/openrc-run
# $Id$


command='/usr/sbin/qubesdb-daemon'
command_args='0'
pidfile='/run/qubes/qubesdb.pid'


depend() {

	need selinux
	need tmpfiles.setup
}
