# app-emulation/qubes-core-qubesdb/files/qubesdb-daemon
#!/sbin/openrc-run
# $Id$


command='/usr/sbin/qubesdb-daemon'
command_args='0'
pidfile='/run/qubes/qubesdb.pid'


depend() {

	need selinux
	need tmpfiles.setup
}

start() {

	ebegin 'Starting qubesdb-daemon'

	# start-stop-daemon hangs =(
	#
	# start-stop-daemon --env --start --exec /usr/sbin/qubesdb-daemon -- 0

	/usr/sbin/qubesdb-daemon 0
}
