# app-emulation/qubes-core-qubesdb/files/qubesdb-daemon_initd
#!/sbin/openrc-run


command='/usr/sbin/qubesdb-daemon'
command_args='0'
pidfile='/run/qubes/qubesdb.pid'

rc_need="selinux"
rc_after="tmpfiles.setup"
