# net-proxy/gpg_socat/files/gpg_socat_initd
#!/sbin/openrc-run

. '/usr/lib/qubes/init/torsocksunix'
