# app-emulation/qubes-gui-agent-linux/files/qubes-gui-agent
#!/sbin/openrc-run
# $Id$


command='/usr/bin/qubes-gui'
command_background='true'
description='QubesOS GUI agent'
pidfile='/run/qubes/qubes-gui.pid'
rc_cgroup_cleanup='yes'
start_stop_daemon_args='--env DISPLAY=:0 --stderr /var/log/qubes/gui-agent.log'

rc_need='qubes-qrexec-agent'
rc_need='qubesdb-daemon'


start_pre() {

	if [ -e '/usr/bin/ck-list-sessions' ]; then {

		rm -f -- '/tmp/qubes-session-env' '/tmp/qubes-session-waiter'

		# start console-kit-daemon

		/usr/bin/ck-list-sessions > /dev/null 2>&1

		# pretend the user is at local console

		touch -- '/var/run/console/user'
	};
	fi
}
