# app-emulation/qubes-core-agent-linux/files/net.qubes
#!/sbin/openrc-run


command='/usr/lib/qubes/setup-ip'
command_background='false'
command_umask='0177'
description='Virtual network interface'
start_stop_daemon_args="--env INTERFACE=\'eth0\' -k 0177"


provide="net"

rc_need="qubesdb-daemon"

rc_after="haveged"
rc_after="iptables"
rc_after="qubes-iptables"
rc_after="qubes-random-seed"
rc_after="urandom"


start() {

	ebegin 'Bringing up virtual network interface'

	INTERFACE="${RC_SVCNAME#net.}" '/usr/lib/qubes/setup-ip'

	eend $?
}

stop() {

	ebegin 'Taking down virtual network interface'

	/bin/ifconfig "${RC_SVCNAME#net.}" down
}
