# app-emulation/qubes-core-agent-linux/files/net.qubes
#!/sbin/openrc-run


depend() {

	need qubesdb-daemon
}

start() {

	ebegin 'Bringing up Qubes virtual network interface'

	INTERFACE='eth0' /usr/lib/qubes/setup-ip
}

stop() {

	ebegin 'Taking down Qubes virtual network interface'

	/bin/ifconfig eth0 down
}
