# app-emulation/qubes-core-agent-linux/files/net.qubes
#!/sbin/openrc-run


depend() {

	need qubesdb-daemon

	after haveged
	after iptables
	after qubes-iptables
	after qubes-random-seed
	after urandom
}

start() {

	ebegin 'Bringing up virtual network interface'

	INTERFACE='eth0' /usr/lib/qubes/setup-ip

	eend $?
}

stop() {

	ebegin 'Taking down virtual network interface'

	/bin/ifconfig eth0 down
}
