# app-emulation/qubes-core-agent-linux/files/net.qubes
#!/sbin/runscript


depend() {

	need qubesdb-daemon
}

start() {

	INTERFACE='eth0' /usr/lib/qubes/setup-ip
}

stop() {

	/bin/ifconfig eth0 down
}
