# app-emulation/qubes-core-agent-linux/files/qubes-firewall
#!/sbin/openrc-run


command='/usr/sbin/qubes-firewall'
command_background='true'
command_umask='0177'
pidfile='/run/qubes/qubes-firewall.pid'
rc_cgroup_cleanup='yes'
required_files='/run/qubes-service/qubes-firewall'
start_stop_daemon_args='-k 0177'

rc_need='qubes-network'
