# app-emulation/qubes-core-agent-linux/files/qubes-random-seed
#!/sbin/openrc-run


rc_keyword='-shutdown'

rc_need='qubesdb-daemon'

rc_after='urandom'

rc_before='haveged'
rc_before='net'


start() {

	# path paranoia
	#
	readonly base64_bin='/usr/bin/base64'
	readonly qubesdb_read_bin='/usr/bin/qubesdb-read'
	readonly qubesdb_rm_bin='/usr/bin/qubesdb-rm'

	ebegin 'Seeding PRNG with entropy from dom0'

	"${qubesdb_read_bin}" /qubes-random-seed | "${base64_bin}" -d > /dev/urandom
	"${qubesdb_rm_bin}" /qubes-random-seed

	eend $?
}

stop() {

	return 0
}
