# app-emulation/qubes-core-agent-linux/files/qubes-qrexec-agent
#!/sbin/openrc-run
# $Id$


command='/usr/lib/qubes/qrexec-agent'
command_background='true'
description='QubesOS RPC agent'
pidfile='/run/qubes/qrexec-agent.pid'
start_stop_daemon_args='--stderr /var/log/qubes/qrexec-agent.log'


depend() {

	# If this isn't a template we don't really need qubes-core.
	#
	after qubes-core
	need selinux
	need tmpfiles.setup
}
