# app-emulation/qubes-core-agent-linux/files/qubes-qrexec-agent
#!/sbin/openrc-run


command='/usr/lib/qubes/qrexec-agent'
command_background='true'
description='QubesOS RPC agent'
pidfile='/run/qubes/qrexec-agent.pid'
start_stop_daemon_args='--stderr /var/log/qubes/qrexec-agent.log'


# If this isn't a template we don't really need qubes-core.
#
rc_after='qubes-core'

rc_need='selinux'
rc_after='tmpfiles.setup'
