# app-emulation/qubes-core-agent-linux/files/qubes-service
#!/sbin/openrc-run


command='/usr/lib/qubes/init/qubes-sysinit.sh'
command_background='false'
command_umask='0177'
description='Qubes services configuration'
start_stop_daemon_args='-k 0177'

rc_keyword='-shutdown'

rc_need='qubesdb-daemon'
rc_need='selinux'
rc_after='tmpfiles.setup'


stop() {

	return 0
}
