# app-emulation/qubes-core-agent-linux/files/qubes-qrexec-agent
#!/sbin/openrc-run
# $Id$

depend() {

	# If this isn't a template we don't really need qubes-core.
	#
	after qubes-core
	need qubes-selinux
	need tmpfiles.setup
}

start() {

	ebegin 'Starting qrexec-agent'

	start-stop-daemon --start --background --exec /usr/lib/qubes/qrexec-agent --make-pidfile --pidfile /run/qubes/qrexec-agent.pid --stderr /var/log/qubes/qrexec-agent.log
}


stop() {

	ebegin 'Stopping qrexec-agent'

	start-stop-daemon --stop --pidfile /run/qubes/qrexec-agent.pid
}
