# app-emulation/qubes-core-agent-linux/files/qubes-service
#!/sbin/openrc-run
# $Id$


command='/usr/lib/qubes/init/qubes-sysinit.sh'
command_background='false'
description='Qubes services configuration'
start_stop_daemon_args='-k 0177'

rc_keyword='-shutdown'

rc_need='qubesdb-daemon'
rc_need='selinux'
rc_need='tmpfiles.setup'


stop() {

	return 0
}
