# app-emulation/qubes-core-agent-linux/files/qubes-network
#!/sbin/openrc-run


command='/usr/lib/qubes/init/network-proxy-setup.sh'
command_backgound='false'
command_umask='0177'
start_stop_daemon_args='-k 0177'
required_files='/run/qubes-service/qubes-network'

rc_keyword='-shutdown'

rc_need='qubes-iptables'


stop() {

	return 0
}
