# app-emulation/qubes-core-agent-linux/files/qubes-netwatcher
#!/sbin/openrc-run


command='/usr/sbin/qubes-netwatcher'
command_background='true'
command_umask='0177'
pidfile='/run/qubes/qubes-netwatcher.pid'
rc_cgroup_cleanup='yes'
required_files='/run/qubes-service/qubes-netwatcher'
start_stop_daemon_args='-k 0177'


rc_need='qubes-firewall'
